# systemctl/Every3hours.service
[Unit]
Description="Service called every three hours"

[Service]
ExecStart=/usr/bin/Every3hours.sh
EnvironmentFile=/etc/scripts.env
