# systemctl/Every3hours.service
[Unit]
Description="Service called every three hours"

[Service]
ExecStart=/usr/bin/Every2hours.sh
