# systemctl/OnBoot.timer
[Unit]
Description=Runs after boot.

[Timer]
OnBootSec=5m
Unit=OnBoot.service

[Install]
WantedBy=timers.target
