# systemctl/OnBoot.service
[Unit]
Description="Service called 15 minutes after pc boot"

[Service]
ExecStart=/usr/bin/OnBoot.sh
