# systemctl/OnBoot.timer
[Unit]
Description=Runs 15 minutes after boot.

[Timer]
OnBootSec=15m
Unit=OnBoot.service

[Install]
WantedBy=timers.target
