# systemctl/OnHourFromBoot.service
[Unit]
Description="Service called every 1h after boot"

[Service]
ExecStart=/usr/bin/OnHourFromBoot.sh
EnvironmentFile=/etc/scripts.env
