# systemctl/OnWorkEnding.service
[Unit]
Description="Service called 15 minutes before work end"

[Service]
ExecStart=/usr/bin/OnWorkEnding.sh
