# systemctl/timer-3hours.timer
[Unit]
Description=Run everyday every 3 hours

[Timer]
OnBootSec=15m
OnCalendar=00/3:00
Unit=Every3hours.service

[Install]
WantedBy=timers.target
