# packages/ubuntu-packages/ppp-server/etc/systemd/system/ppp-server.service
[Unit]
Description="PPP server"

[Service]
ExecStart=/usr/sbin/pppd /dev/ttyECRPPP
#EnvironmentFile=/etc/scripts.env
Restart=on-failure
RestartSec=15s

[Install]
WantedBy=multi-user.target
